# local/etc/profile.d/path_adm_bin.sh
#! /bin/false
# Inserts bin directory for additional system-wide binaries into $PATH.

set -- \
	--append /home/group/admdevel/bin \
	--append /home/project/admdevel/bin \
	--stop "$@"
. /usr/local/libexec/xworld/sh/path_hyec3v5m8kd1vjs8k7d1wce62.sh
