# local/etc/profile.d/shell-goodies.sh
#! /bin/false
# Add support for extensible bash- and shell-plugins.
# Despite the name, should work for most Bourne-derived
# shells, not just for bash.
. /usr/local/etc/site/bash-goodies
